# services/pulpmonad.service
[Unit]
Description=PulpMonad

[Service]
Type=simple
StandardOutput=journal
StandardError=journal
ExecStart=${XMONAD_CACHE_DIR}/${XMONAD_NAME}

[Install]
